# pkg/deck-tenant/usr/lib/systemd/user/deck-tenant-login.path
[Unit]
Description=deck-tenant: watch for Steam account changes

[Path]
PathModified=%h/.steam/steam/config/loginusers.vdf

[Install]
WantedBy=default.target
