# pkg/deck-tenant/usr/lib/systemd/user/deck-tenant-guard.service
[Unit]
Description=deck-tenant: close/reopen registered apps on Steam session and account boundaries

[Service]
Type=simple
ExecStart=%h/.local/share/deck-pkgs/usr/bin/deck-tenant-guard
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
