# dist/deck-tenant-login.service
[Unit]
Description=deck-tenant: handle a Steam account change

[Service]
Type=oneshot
ExecStart=/usr/bin/deck-tenant-guard handle-login
